# linux-amd64.service
[Unit]
Description=Uptime Monitoring service
After=network.target
Wants=network-online.target

[Service]
Restart=always
Type=simple
ExecStart=/srv/monitor/heartbeat-sender-0.1.0-linux-x86_64
WorkingDirectory=/srv/monitor

[Install]
WantedBy=multi-user.target
